# tests/container/reboot/bootc-finish-test-reboot.service
[Unit]
ConditionPathExists=!/etc/initrd-release 
After=local-fs.target
RequiresMountsFor=/run/bootc-test-reboot
Before=bootc-test-reboot.service
PartOf=bootc-test-reboot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStop=touch /run/bootc-test-reboot/success
